# src/offspot_demo/systemd-unit/demo-watcher.service
[Unit]
Description=demo-watcher

[Service]
User=root
ExecStart=/root/env/bin/demo-watcher
